# set-hyprsunset.timer
[Unit]
Description=Run display cycle every 10 minutes
Requires=display-cycle-user.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
